# backend/systemd/dimaist.service
[Unit]
Description=Dimaist Task Management API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/dimaist
ExecStart=/root/dimaist/dimaist
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

# Environment variables (if needed)
# Environment="PORT=3000"
# Environment="DATABASE_URL=postgres://..."
EnvironmentFile=/root/dimaist/.env

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
